# roles/server/files/docker/container_maintenance.timer
[Unit]
Description=Docker container maintenance timer

[Timer]
OnCalendar=*-*-* 5:00:00
Persist=true

[Install]
WantedBy=timers.target